# systemd/file-classifier.path
[Unit]
Description=File classifier for downloads
[Path]
PathChanged={{HOME_DIR}}/Downloads
[Install]
WantedBy=default.target
